# x11vnc/supervisord.conf
; supervisor config file Xvnc4

[unix_http_server]
file=%(ENV_HOME)s/.supervisor/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[supervisord]
logfile=%(ENV_HOME)s/.supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=%(ENV_HOME)s/.supervisor/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=%(ENV_HOME)s/.supervisor/log            ; ('AUTO' child log dir, default $TEMP)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
nocleanup=false              ; (don't clean up tempfiles at start;default false)
;umask=022                   ; (process file creation umask;default 022)
;user=chrism                 ; (default is current user, required if root)
;identifier=supervisor       ; (supervisord identifier, default is 'supervisor')
;directory=%(ENV_HOME)s          ; (default is not to cd during start)
;environment=KEY="value"     ; (key value pairs to add to environment)
;strip_ansi=false            ; (strip ansi escape codes in logs; def. false)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///%(ENV_HOME)s/.supervisor/supervisor.sock ; use a unix:// URL  for a unix socket

[program:x11vnc]
command=x11vnc -env FD_OPTS="-nolisten tcp -c r" -env FD_GEOM="1024x768x24" -display WAIT:cmd=FINDCREATEDISPLAY-Xvfb -no6 -noipv6 -noxdamage -noxfixes -noxrecord -rfbport 5920 -ping 5 -localhost -repeat -forever -nevershared -nopw
; the program (relative uses PATH, can take args)
numprocs=1                     ; number of processes copies to start (def 1)
directory=/                    ; directory to cwd to before exec (def no cwd)
priority=1                     ; the relative start priority (default 999)
autostart=true                 ; start at supervisord start (default: true)
autorestart=true               ; when to restart if exited after running (def: unexpected)
redirect_stderr=true           ; redirect proc stderr to stdout (default false)
stdout_logfile=%(ENV_HOME)s/.supervisor/xvnc.log ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=1MB    ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10      ; # of stdout logfile backups (default 10)
startsecs=5                    ; # of secs prog must stay up to be running (def. 1)
;process_name=%(program_name)s ; process_name expr (default %(program_name)s)
;umask=022                     ; umask for process (default None)
;startretries=3                ; max # of serial start failures when starting (default 3)
;exitcodes=0,2                 ; 'expected' exit codes used with autorestart (default 0,2)
;stopsignal=QUIT               ; signal used to kill process (default TERM)
;stopwaitsecs=10               ; max num secs to wait b4 SIGKILL (default 10)
;stopasgroup=false             ; send stop signal to the UNIX process group (default false)
;killasgroup=false             ; SIGKILL the UNIX process group (def false)
;user=chrism                   ; setuid to this UNIX account to run the program
;stdout_capture_maxbytes=1MB   ; number of bytes in 'capturemode' (default 0)
;stdout_events_enabled=false   ; emit events on stdout writes (default false)
;stderr_logfile=/a/path        ; stderr log path, NONE for none; default AUTO
;stderr_logfile_maxbytes=1MB   ; max # logfile bytes b4 rotation (default 50MB)
;stderr_logfile_backups=10     ; # of stderr logfile backups (default 10)
;stderr_capture_maxbytes=1MB   ; number of bytes in 'capturemode' (default 0)
;stderr_events_enabled=false   ; emit events on stderr writes (default false)
;environment=A="1",B="2"       ; process environment additions (def no adds)
;serverurl=AUTO                ; override serverurl computation (childutils)

[program:novnc]
command=/usr/share/novnc/utils/launch.sh --vnc localhost:%(ENV_VNCPORT)s --listen %(ENV_WSPORT)s
; the program (relative uses PATH, can take args)
numprocs=1                     ; number of processes copies to start (def 1)
directory=/                    ; directory to cwd to before exec (def no cwd)
priority=2                     ; the relative start priority (default 999)
autostart=true                 ; start at supervisord start (default: true)
autorestart=true               ; when to restart if exited after running (def: unexpected)
redirect_stderr=true           ; redirect proc stderr to stdout (default false)
stdout_logfile=%(ENV_HOME)s/.supervisor/novnc.log ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=1MB    ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10      ; # of stdout logfile backups (default 10)

[program:fluxbox]
command=/usr/bin/startfluxbox
; the program (relative uses PATH, can take args)
numprocs=1                     ; number of processes copies to start (def 1)
directory=/                    ; directory to cwd to before exec (def no cwd)
priority=3                     ; the relative start priority (default 999)
autostart=true                 ; start at supervisord start (default: true)
autorestart=true               ; when to restart if exited after running (def: unexpected)
redirect_stderr=true           ; redirect proc stderr to stdout (default false)
stdout_logfile=%(ENV_HOME)s/.supervisor/xsession.log ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=1MB    ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10      ; # of stdout logfile backups (default 10)

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

; [include]
; files = /etc/supervisor/conf.d/*.conf %(ENV_HOME)s/.supervisor/conf.d/*.conf
